# config/autostart/wal.desktop
[Desktop Entry]
Exec=wal -R
Icon=
Name=wal
Path=
Terminal=False
Type=Application
